# .config/autostart/davm-docs.desktop
[Desktop Entry]
Name=DAVM Docs
Exec=sleep 2; kstart5 firefox-bin http://localhost/wordpress
Type=Application
Terminal=false
