# .config/autostart/davm-docs.desktop
[Desktop Entry]
Name=DAVM Docs
Exec=kstart5 --maximize firefox-bin http://localhost/wordpress
Type=Application
Terminal=false
