# .config/autostart/davm-launcher.desktop
[Desktop Entry]
Name=DAVM Demo Launcher 
Exec=/local/examples/launcher/serve.sh &
Type=Application
Terminal=false
