# .config/autostart/davm-version-check.desktop
[Desktop Entry]
Name=DAVM Version Check 
Exec=sleep 5; kstart5 /local/bin/checkVersion
Type=Application
Terminal=false
